# OS/airootfs/etc/systemd/system/setup.service
[Unit]
Description=Run setup python script

[Service]
ExecStart=/bin/python /root/.setup

[Install]
WantedBy=multi-user.target
