# OS/airootfs/etc/systemd/system/multi-user.target.wants/defaultwall.service
[Unit]
Description=Change default wallpaper for cscOS

[Service]
ExecStart=/root/.wallpaper/defaultwall.sh

[Install]
WantedBy=multi-user.target
